# scraper/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

0 2 * * * root cd /app/scraper && /app/scraper/run_scraper.sh >> /var/log/scraper.log 2>&1
